# service/sipwise-backup.service
[Unit]
Description=Sipwise Backup CLI Service
After=network.target
Documentation=https://github.com/TimothySteinert/sipwise-backup

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/sipwise-backup/scheduler.py
WorkingDirectory=/opt/sipwise-backup
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
